# files/etc/services.d/crond/run
#!/usr/bin/execlineb -P

crond -f -l 8
#-L /var/log/crond.log
